# multiappproxy/rootfs/etc/services.d/api/finish
#!/command/with-contenv bashio
# ==============================================================================
# API server stopped
# ==============================================================================

bashio::log.warning "API server stopped, restarting..."
