# multiappproxy/rootfs/etc/services.d/nginx/finish
#!/command/with-contenv bashio
# ==============================================================================
# Nginx stopped
# ==============================================================================

bashio::log.warning "Nginx stopped, restarting..."
