# multiappproxy/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Start nginx service
# ==============================================================================

bashio::log.info "Starting nginx..."

# Démarrer nginx en arrière-plan
nginx -g "daemon off;" &
NGINX_PID=$!

# Suivre les logs
tail -f /var/log/nginx/access.log /var/log/nginx/error.log &

# Attendre nginx
wait $NGINX_PID
