# multiappproxy/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Start nginx service
# ==============================================================================

bashio::log.info "Starting nginx..."

# Run nginx in the foreground (s6 manages the process)
nginx -g "daemon off;" &
NGINX_PID=$!

# Tail logs so they appear in the s6 journal
tail -f /var/log/nginx/access.log /var/log/nginx/error.log &

# Wait for nginx to exit
wait $NGINX_PID
